# satis/supervisor/conf.d/1-cron.conf
[program:cron]
command=/usr/sbin/crond -f
autostart = true
autorestart=true
priority=5