# satis/supervisor/conf.d/2-atd.conf
[program:atd]
command=/usr/sbin/atd -f
autostart = true
autorestart=true
priority=5